# service/gateone.conf
description "GateOne Service"
author "Jarod Law <dingyong.law@gmail.com>"

start on runlevel [234]
stop on runlevel [0156]

chdir /opt/gateone
exec /opt/gateone/gateone.py --disable_ssl
respawn
